# phpsyslogng/files/phpsyslogng.cron
01 0 * * * root /usr/bin/php /var/www/php-syslog-ng/scripts/logrotate.php
